# tools/crontab
0 21 * * 1-6 /usr/sml/MLW/tools/daily-wrapper
0 21 * * 0 /usr/sml/MLW/tools/weekly-wrapper
